# px-deploy.cron
0 8,16 * * * bash $HOME/.px-deploy/px-deploy_cron.sh
