# px-deploy_mac.cron
0 8,16 * * * bash $HOME/.px-deploy/px-deploy_mac_cron.sh
